# nginx/configs/etc/systemd/system/tarea8.socket
[Unit]
Description=tarea8 socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target